# root/etc/s6-overlay/s6-rc.d/sonarr/run
#!/usr/bin/execlineb -P

mono --debug /opt/sonarr/bin/Sonarr.exe --nobrowser --data=/container/data
